# wg-bot.service
[Unit]
Description=WireGuard Telegram Bot
After=network.target

[Service]
ExecStart=/home/USERNAME/scripts/wg-bot/wg-bot -vv
WorkingDirectory=/home/USERNAME/scripts/wg-bot
Environment=TELEGRAM_TOKEN=your_token_here
Environment=WG_INTERFACE=wg0
Environment=ALLOWED_USER=123456789
Restart=on-failure
RestartSec=5
User=USERNAME
Group=USERNAME
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=wg-bot

[Install]
WantedBy=multi-user.target
